# deploy/kamelia.service
[Unit]
Description=Kamelia-Med FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/kamelia/backend
Environment="PATH=/opt/kamelia/backend/.venv/bin"
Environment="KAMELIA_SECRET=CHANGE_ME_IN_PRODUCTION"
ExecStart=/opt/kamelia/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
